# infra/cron/crontab
MAILTO=""
* * * * * root /usr/local/bin/php /app/bin/console schedule:run >> /var/log/cron.log 2>&1
